# initscript/gentoo/avahi-dnsconfd.in
#!/sbin/runscript
# $Id$
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/www/www.gentoo.org/raw_cvs/gentoo-x86/net-dns/avahi/files/avahi-dnsconfd.initd,v 1.1 2005/08/22 03:47:19 swegener Exp $

opts="${opts} reload"

depend() {
	need avahi-daemon
	use net
}

start() {
	ebegin "Starting avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -D
	eend $?
}

stop() {
	ebegin "Stopping avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -k
	eend $?
}

reload() {
	ebegin "Reloading avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -r
	eend $?
}
