# initscript/Gentoo/avahi-daemon
#!/sbin/runscript
# Avahi ZeroConf Daemon init script
# Many thanks to the Avahi developer's and Sebastien Estienne

DAEMON=/usr/bin/avahi-daemon

depend () {
	need dbus
}

start() {
	ebegin "Starting avahi"
        $DAEMON -D
	eend $?
}

stop() {
	ebegin "Stopping avahi"
	$DAEMON -k || /bin/true
	eend $?
}

reload() {
	ebegin "Reloading avahi configuration"
        $DAEMON -r
	eend $?
}

